# crontab
#!/bin/bash
# min   hour    day     month   weekday command
* * * * * python3 /home/ubuntu/app/live-now-job/hello.py  >> /var/log/cron/cron-`date +\%Y-\%m-\%d`.log 2>&1