# crontab
#!/bin/bash
# min   hour    day     month   weekday command
1 * * * * python3 /app/live-now-job/hello.py  >> /var/log/cron.log 2>&1